# build/lib/fair_ros/watchdog/fair-ros-watchdog.service
[Unit]
Description=fair-ros recording assistant (mission context watchdog)
Documentation=https://github.com/example/fair-ros
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m fair_ros.watchdog.watchdog
Restart=always
RestartSec=5
# The watchdog must keep working when ROS 2 / Docker are down; it only needs
# the spool. Run as root group fair-ros so operators can read the spool.
UMask=0002
Nice=10

[Install]
WantedBy=multi-user.target
